# tools/ubuntu_startup/rc.local
#!/bin/sh
echo "rc,local start success" > /usr/local/test.log
ifconfig enp0s3 up
dhclient -v
exit 0
